# deepseekbot.service
[Unit]
Description=My Python Script
After=network.target

[Service]
ExecStart=/usr/bin/python3.12 main.py
WorkingDirectory=/root/Personal_AI_TG_bot
User=root
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
